# fan_control.service
[Unit]
Description=PWM Fan Control Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/pwm-fan-control.py
Restart=on-failure
RestartSec=5s
User=pi
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
