# fan_control.service
[Unit]
Description=PWM Fan Control Service
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/rpi-fan-control/pwm-fan-control.py
Restart=on-failure
StandardOutput=append:/home/pi/rpi-fan-control/fan_log.txt
StandardError=inherit

[Install]
WantedBy=multi-user.target
